# config/autostart/sleep.desktop
[Desktop Entry]
Type=Application
Exec=xchat
Hidden=false
X-MATE-Autostart-enabled=true
Name[en_US]=XChat
Name=XChat
Comment[en_US]=
Comment=
